# docs/docker/docker copy.slice
[Unit]
Description=test.slice
DefaultDependencies=no
Before=slices.target

[Install]
WantedBy=slices.target

[Slice]
