# docs/docker/docker.slice
[Unit]
Description=docker.slice
Before=slices.target

# [Install]
# WantedBy=slices.target

[Slice]
